# roles/backup/templates/restic-backup-systemd.service.j2
[Unit]
Description=Restic backup for {{ inventory_hostname }}
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/restic-backup
StandardOutput=journal
StandardError=journal
