# roles/backup/templates/restic-backup-systemd.timer.j2
[Unit]
Description=Daily restic backup timer for {{ inventory_hostname }}

[Timer]
OnCalendar=*-*-* {{ backup_schedule }}
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
